# backend/tip-generator-worker.service
[Unit]
Description=TIP Generator Celery Worker
After=network.target postgresql.service redis.service tip-generator.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/tip-generator/backend
Environment="PATH=/opt/tip-generator/backend/venv/bin"
Environment="PYTHONPATH=/opt/tip-generator/backend"
EnvironmentFile=/opt/tip-generator/backend/.env
ExecStart=/opt/tip-generator/backend/venv/bin/celery -A celery_app worker --loglevel=info --pool=threads --concurrency=3 -n tip-worker@%%h
KillSignal=SIGKILL
TimeoutStopSec=600
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
